# deploy/cloudru/audiogid-api.service
[Unit]
Description=Audiogid API Service
After=network.target postgresql.service

[Service]
Type=simple
User=user1
Group=user1
WorkingDirectory=/opt/audiogid/api
Environment="PATH=/home/user1/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/audiogid/api/.env
ExecStart=/home/user1/.local/bin/uvicorn index:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
